# root/etc/update-motd.d/05-banner
#!/bin/sh
cat <<'EOF'
       _____ _           _       _ _ _ _____ _____
      |  |  | |_ _ _ ___| |_ _ _| | | | __  |_   _|
      |  |  | . | | |   |  _| | | | | |    -| | |
      |_____|___|___|_|_|_| |___|_____|__|__| |_|

UbuntuWRT 20.02 - https://github.com/cilix-lab/ubuntu-wrt
EOF
