# rootfs/etc/network/if-pre-up.d/iptables-restore
#!/bin/sh

if [ "$IFACE" = "wan" ]; then
	iptables-restore < /etc/iptables.rules
	ip6tables-restore < /etc/ip6tables.rules
fi

exit 0
